# deploy/brilliant-wifi-watchdog.service
[Unit]
Description=Brilliant Wi-Fi connectivity watchdog
After=network.target connman.service
Wants=network.target

[Service]
Type=simple
EnvironmentFile=-/etc/brilliant-mqtt.env
ExecStart=/usr/bin/python3.10 -m brilliant_wifi_watchdog.run
WorkingDirectory=/var/brilliant-mqtt/wifi_watchdog
Environment=PYTHONPATH=/var/brilliant-mqtt/wifi_watchdog
Restart=always
RestartSec=10
Nice=10
MemoryMax=48M
CPUQuota=10%

[Install]
WantedBy=multi-user.target
